# scripts/com.eran.mundial.ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.eran.mundial.ingest</string>
  <key>ProgramArguments</key>
  <array>
    <string>__UV__</string>
    <string>run</string><string>python</string><string>-m</string>
    <string>ingestion.run_daily</string>
  </array>
  <key>WorkingDirectory</key><string>__PROJECT__</string>
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
  <key>StandardOutPath</key><string>__PROJECT__/data/ingest.log</string>
  <key>StandardErrorPath</key><string>__PROJECT__/data/ingest.err</string>
</dict>
</plist>
